# deploy/polymarket-collector.service
[Unit]
Description=Polymarket Esports Data Collector
After=network.target

[Service]
Type=simple
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-collector
ExecStart=/opt/polymarket-collector/venv/bin/python main.py --continuous
Restart=always
RestartSec=10

# Ensure logs appear immediately in journald
Environment=PYTHONUNBUFFERED=1

# Trading credentials and configuration (created by setup.sh)
EnvironmentFile=-/etc/polymarket/trading.env

# Memory guard — bumped 600M → 800M after periodic trades polling was added
# (Phase 2). Box has 1.9G total; api service is capped at 256M, leaves slack
# for kernel + nginx. OOM-restart was masking dropped trades capture.
MemoryMax=800M

# Logging via journald (view with: journalctl -u polymarket-collector -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polymarket-collector

[Install]
WantedBy=multi-user.target
